# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:api]
command=node /app/api/dist/index.mjs
directory=/app/api/dist
autostart=true
autorestart=true
priority=20
; PORT is required by index.ts — set to 3000 (internal only, nginx proxies in front)
; DATABASE_URL is injected from Azure Key Vault via Container App env vars
environment=PORT="3000",NODE_ENV="production",DATABASE_URL="%(ENV_DATABASE_URL)s"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
